# detectcard.service
[Unit]
Description=Detect card script for dash cam import.

[Service]
ExecStart=/usr/local/bin/detectcard.py